# containers/go-1.23/supervisord.conf
[supervisord]
nodaemon=true

[program:bridge]
command=deno run --allow-net /usr/bridge/src/index.ts
directory=/usr/bridge
autostart=true
autorestart=true
stderr_logfile=/var/log/deno.err.log
stdout_logfile=/var/log/deno.out.log

[program:nginx]
command=nginx -g 'daemon off;'
autostart=true
autorestart=true
stderr_logfile=/var/log/nginx.err.log
stdout_logfile=/var/log/nginx.out.log
